# x11VNC/autostart/1.conf
start on login-session-start

script

sudo <PATH>/config.sh start <PORT>

end script
