# x11VNC/autostart/2.conf
description "x11vnc"

start on runlevel [2345]
stop on runlevel [^2345]

console log

respawn
respawn limit 20 5

exec <PATH>/config.sh start <PORT>
